# mqttcamera.service
[Unit]
Description=MQTT Camera

[Service]
ExecStart=/usr/local/bin/mqttcamera.sh
Restart=on-abort

[Install]
WantedBy=multi-user.target
